# vpn/connman-vpn.service.in
[Unit]
Description=ConnMan VPN service
Requires=dbus.socket
After=dbus.socket

[Service]
Type=dbus
BusName=net.connman.vpn
SmackProcessLabel=System
ExecStart=@sbindir@/connman-vpnd -n
StandardOutput=null
CapabilityBoundingSet=~CAP_MAC_ADMIN
CapabilityBoundingSet=~CAP_MAC_OVERRIDE

[Install]
WantedBy=multi-user.target
