# vpn/connman-vpn.service.in
[Unit]
Description=ConnMan VPN service
Requires=dbus.socket
After=dbus.socket

[Service]
Type=dbus
User=network_fw
Group=network_fw
BusName=net.connman.vpn
SmackProcessLabel=System
ExecStart=@bindir@/connman-vpnd -n
StandardOutput=null
Capabilities=cap_setgid,cap_net_admin,cap_net_bind_service,cap_net_broadcast,cap_net_raw,cap_dac_override=i
SecureBits=keep-caps

[Install]
WantedBy=multi-user.target
