# src/connman.service.in
[Unit]
Description=Connection service
After=net-config.service
DefaultDependencies=no

[Service]
Type=dbus
BusName=net.connman
Restart=on-failure
SmackProcessLabel=System
ExecStart=@sbindir@/connmand -n --noplugin vpn
User=network_fw
Group=network_fw
StandardOutput=null

[Install]
WantedBy=multi-user.target
